# natssl-master.service
[Unit]
Description=NATSSL Master (Private CA)
Documentation=https://github.com/iskyneon/natssl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/natssl --mode=master --config=/etc/natssl/config.yaml
Restart=on-failure
RestartSec=3s

# The master binds :443; it does NOT use raw sockets (no ICMP on this path),
# so CAP_NET_BIND_SERVICE is the only privilege required.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

# State directory (data_dir). systemd creates /var/lib/natssl (0750) for us.
StateDirectory=natssl
StateDirectoryMode=0750

# Filesystem hardening: only the data dir is writable.
ProtectSystem=strict
ReadWritePaths=/var/lib/natssl
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

# Networking surface.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
