# natssl-client.service
[Unit]
Description=NATSSL Client (Cert Store)
Documentation=https://github.com/iskyneon/natssl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/natssl --mode=client --config=/etc/natssl/config.yaml
Restart=on-failure
RestartSec=3s

# The client listens on :8443 (>1024, no bind privilege needed). CAP_NET_RAW is
# kept as an ambient capability so the disaster-recovery promotion path
# (`ping` during the safety chain) works even under NoNewPrivileges=true.
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
NoNewPrivileges=true

# State directory (data_dir).
StateDirectory=natssl
StateDirectoryMode=0750

# NOTE: the client installs the Root CA into the OS trust store, so it must be
# able to write to the system anchor directories and run update-ca-* tools.
# ProtectSystem is therefore relaxed (full, not strict) and the trust-store
# paths are explicitly allowed.
ProtectSystem=full
ReadWritePaths=/var/lib/natssl /usr/local/share/ca-certificates /etc/pki/ca-trust/source/anchors /etc/ssl/certs
PrivateTmp=true
ProtectKernelModules=true
RestrictSUIDSGID=false
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

[Install]
WantedBy=multi-user.target
